# supervisor.ini
[program:gin-hong-api-queue]
directory=/golang/go/bin
command=/golang/go/bin/queue -env pro
autostart=true
autorestart=true
startsecs=1
user = root
stderr_logfile=/tmp/queue.log
stdout_logfile=/tmp/queue.log
redirect_stderr = true
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 20